# .quadlet/beprepared-dev.container
# bePrepared development container
# Install to: ~/.config/containers/systemd/beprepared-dev.container

[Unit]
Description=bePrepared development container

[Container]
Image=localhost/beprepared-dev:latest
ContainerName=beprepared-dev
WorkingDir=/workspace
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
