# deploy/quadlet/beprepared-db.container
# bePrepared — MariaDB container
# Quadlet file: ~/.config/containers/systemd/beprepared-db.container

[Unit]
Description=bePrepared MariaDB database
After=beprepared-pod.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=docker.io/library/mariadb:11
Pod=beprepared.pod
ContainerName=beprepared-db

EnvironmentFile=%h/0_opencode/bePrepared/.env

Volume=beprepared-db.volume:/var/lib/mysql:Z

# Functional healthcheck: the official mariadb:11 image ships healthcheck.sh
# which validates connections and InnoDB initialisation.
# HealthOnFailure=kill lets systemd own the restart so StartLimitBurst applies.
HealthCmd=healthcheck.sh --connect --innodb_initialized
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=60s
HealthRetries=3
HealthOnFailure=kill

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=beprepared-pod.service
