# deploy/quadlet/beprepared-api.container
# bePrepared — API container
# Quadlet file: ~/.config/containers/systemd/beprepared-api.container

[Unit]
Description=bePrepared REST API (Bun + Elysia)
After=beprepared-db.container

[Container]
Image=localhost/beprepared-api:latest
Pod=beprepared.pod
ContainerName=beprepared-api

EnvironmentFile=%h/0_opencode/bePrepared/deploy/.env
PodmanArgs=--read-only
PodmanArgs=--tmpfs=/tmp:rw,nosuid,nodev,noexec,size=64m
PodmanArgs=--security-opt=no-new-privileges
PodmanArgs=--cap-drop=all

# Automatically restart the container when its HEALTHCHECK transitions to
# "unhealthy" (i.e. 3 consecutive /health probe failures).

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=beprepared-pod.service
