# deploy/quadlet/beprepared-db.container
# bePrepared — MariaDB container
# Quadlet file: ~/.config/containers/systemd/beprepared-db.container

[Unit]
Description=bePrepared MariaDB database
After=beprepared-pod.service

[Container]
Image=docker.io/library/mariadb:11
Pod=beprepared.pod
ContainerName=beprepared-db

EnvironmentFile=%h/0_opencode/bePrepared/deploy/.env

Volume=beprepared-db:/var/lib/mysql:Z

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=beprepared-pod.service
