# deploy/quadlet/beprepared-worker.container
# bePrepared — Worker container (hourly alert scheduler)
# Quadlet file: ~/.config/containers/systemd/beprepared-worker.container

[Unit]
Description=bePrepared background worker (alert scheduler)
After=beprepared-api.container
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=localhost/beprepared-worker:latest
Pod=beprepared.pod
ContainerName=beprepared-worker

EnvironmentFile=%h/0_opencode/bePrepared/deploy/.env
PodmanArgs=--read-only
PodmanArgs=--tmpfs=/tmp:rw,nosuid,nodev,noexec,size=64m
PodmanArgs=--security-opt=no-new-privileges
PodmanArgs=--cap-drop=all

# Functional healthcheck: /tmp/worker.ready must exist and be ≤30 min old
# (2× the 15-min run interval). Stale = worker silently stopped scheduling.
# HealthOnFailure=kill lets systemd own the restart so StartLimitBurst applies.
HealthCmd=find /tmp/worker.ready -mmin -30 | grep -q .
HealthInterval=30s
HealthTimeout=5s
HealthStartPeriod=120s
HealthRetries=3
HealthOnFailure=kill

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=beprepared-pod.service
