# deploy/quadlet/beprepared-frontend.container
# bePrepared — Frontend container (Next.js)
# Quadlet file: ~/.config/containers/systemd/beprepared-frontend.container

[Unit]
Description=bePrepared frontend (Next.js)
After=beprepared-api.container

[Container]
Image=localhost/beprepared-frontend:latest
Pod=beprepared.pod
ContainerName=beprepared-frontend

EnvironmentFile=%h/bePrepared/deploy/.env
Environment=PORT=9999

[Service]
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=beprepared-pod.service
