# deploy/quadlet/beprepared-api.container
# bePrepared — API container
# Quadlet file: ~/.config/containers/systemd/beprepared-api.container

[Unit]
Description=bePrepared REST API (Bun + Elysia)
After=beprepared-db.container
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=localhost/beprepared-api:latest
Pod=beprepared.pod
ContainerName=beprepared-api

EnvironmentFile=%h/0_opencode/bePrepared/.env
PodmanArgs=--read-only
PodmanArgs=--tmpfs=/tmp:rw,nosuid,nodev,noexec,size=64m
PodmanArgs=--security-opt=no-new-privileges
PodmanArgs=--cap-drop=all

# Functional healthcheck: hits /health which validates DB connectivity.
# HealthOnFailure=kill lets systemd own the restart so StartLimitBurst applies.
HealthCmd=wget -qO- http://127.0.0.1:3001/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=60s
HealthRetries=3
HealthOnFailure=kill

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=beprepared-pod.service
