# deploy/quadlet/beprepared-worker.container
# bePrepared — Worker container (hourly alert scheduler)
# Quadlet file: ~/.config/containers/systemd/beprepared-worker.container

[Unit]
Description=bePrepared background worker (alert scheduler)
After=beprepared-api.container

[Container]
Image=localhost/beprepared-worker:latest
Pod=beprepared.pod
ContainerName=beprepared-worker

EnvironmentFile=%h/0_opencode/bePrepared/deploy/.env
PodmanArgs=--read-only
PodmanArgs=--tmpfs=/tmp:rw,nosuid,nodev,noexec,size=64m
PodmanArgs=--security-opt=no-new-privileges
PodmanArgs=--cap-drop=all

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=beprepared-pod.service
